# 7.Script/ubuntu/overlay/etc/update-motd.d/00-header
#!/bin/sh
#

[ -r /etc/lsb-release ] && . /etc/lsb-release

printf "\n"
echo "\e[1;34m 
                                                  
                                                  
                    .:-==.                        
                 :+***##***:                      
               .+#*****-=*#=                      
              -*****#=.+**#-                      
             +#**##*- +#*#*.                      
             ==---=: .*+==:::                     
           .:-=======-:.  -=-.:-. -:              
        .-+****++******+=:  ..:-=+**+-            
       -****+:::::-*******+. -*********:          
      +***+=.-==== =********: :*********          
     -*+**+ ===-:=- +********: -********.         
     =****..===:.=::*********=  .::----:          
     -****=:.. :-:.**********+ :==+++++++-:       
     .*******+-::-+**********= ************+-     
      :**********************. ...::-==++****+:   
       :+*******************:  ====----==+**+*=   
         -+**************+-    :+*************-   
           .-=++****++=-.        :=+*******+=-.   
                 ..                 .::::.        
                                                                                                 
  ____   _               _                                ____   _ 
 | __ ) | | _   _   ___ | |__    ___  _ __  _ __  _   _  |  _ \ (_)
 |  _ \ | || | | | / _ \| '_ \  / _ \| '__|| '__|| | | | | |_) || |
 | |_) || || |_| ||  __/| |_) ||  __/| |   | |   | |_| | |  __/ | |
 |____/ |_| \__,_| \___||_.__/  \___||_|   |_|    \__, | |_|    |_|
                                                  |___/            
\e[0m"

if [ -z "$DISTRIB_DESCRIPTION" ] && [ -x /usr/bin/lsb_release ]; then
	# Fall back to using the very slow lsb_release utility
	DISTRIB_DESCRIPTION=$(lsb_release -s -d)
fi
printf "\n"
printf "Welcome to %s (%s %s %s)\n" "$DISTRIB_DESCRIPTION" "$(uname -o)" "$(uname -r)" "$(uname -m)"
