# 7.Script/ubuntu/overlay/etc/profile.d/x11.sh
export DISPLAY=${DISPLAY:-:0}
export XSERVER_FREEZE_DISPLAY=/tmp/.freeze_xserver
